# supervisor/templates/default/program.conf.erb
#
# /etc/supervisor/conf.d/<%= @name %>.conf
#
# Generated by Chef for <%= node[:fqdn] %>
#

[program:<%= @name %>]
command=<%= @command %>
stdout_logfile=<%= @stdout_logfile %>
autostart=<%= @autostart %>
autorestart=<%= @autorestart %>
user=<%= @user %>
startsecs=<%= @startsecs %>
numprocs=<%= @numprocs %>